# cron/slurm-node-health.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Run every 5 minutes as the slurm user to detect and resume stuck nodes.
# Install: sudo bash scripts/node-health/install.sh
*/5 * * * * slurm /usr/local/sbin/slurm-node-health
