# cron/slurm-node-health.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Run every 5 minutes as the slurm user to detect and resume stuck nodes.
# Install: cp cron/slurm-node-health.cron /etc/cron.d/slurm-node-health
*/5 * * * * slurm /gpfs/mskmind_ess/limr/repos/slurm-devops/scripts/node-health/fix-stuck-nodes.sh
